# supervisord.conf
[program:dropwizard]
command = java -jar /opt/AppOne.jar server
[supervisord]
nodaemon = true
